# systemd/superbot.service
[Unit]
Description=SuperBot Slack Bridge
Documentation=https://github.com/irismed-oso/super_bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bot
Group=bot
WorkingDirectory=/home/bot/super_bot
EnvironmentFile=/home/bot/.env
ExecStart=/home/bot/super_bot/.venv/bin/python -m bot.app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=superbot

# Resource limits
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
